# certmate.service
[Unit]
Description=CertMate SSL Certificate Manager
Documentation=https://github.com/fabriziosalmi/certmate
After=network.target network-online.target
Wants=network-online.target

[Service]
Type=simple
User=certmate
Group=certmate
WorkingDirectory=/opt/certmate
Environment=PATH=/opt/certmate/venv/bin
Environment=API_BEARER_TOKEN=change-this-secure-token
# Single worker + threads, matching the Dockerfile (#411). Sessions, session
# revocation and the login rate limiter are per-process in-memory state: with
# 4 workers, revoking a compromised admin's session only took effect on the
# worker that handled the request (so ~3 requests in 4 still validated), the
# login rate limit was multiplied by the worker count, and APScheduler ran a
# duplicate renewal job per worker.
ExecStart=/opt/certmate/venv/bin/gunicorn --bind 0.0.0.0:8000 --workers 1 --threads 8 app:app
ExecReload=/bin/kill -s HUP $MAINPID
Restart=always
RestartSec=10
KillMode=mixed
TimeoutStopSec=5

# Security settings
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/certmate/certificates /opt/certmate/data
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true

[Install]
WantedBy=multi-user.target